# contrib/init/nestcoind.conf
description "Nestcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env NESTCOIND_BIN="/usr/bin/nestcoind"
env NESTCOIND_USER="nestcoin"
env NESTCOIND_GROUP="nestcoin"
env NESTCOIND_PIDDIR="/var/run/nestcoind"
# upstart can't handle variables constructed with other variables
env NESTCOIND_PIDFILE="/var/run/nestcoind/nestcoind.pid"
env NESTCOIND_CONFIGFILE="/etc/nestcoin/nestcoin.conf"
env NESTCOIND_DATADIR="/var/lib/nestcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # nestcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$NESTCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run nestcoind."
        echo "The setting must appear in $NESTCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Nestcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$NESTCOIND_PIDDIR"
    chmod 0755 "$NESTCOIND_PIDDIR"
    chown $NESTCOIND_USER:$NESTCOIND_GROUP "$NESTCOIND_PIDDIR"
    chown $NESTCOIND_USER:$NESTCOIND_GROUP "$NESTCOIND_CONFIGFILE"
    chmod 0660 "$NESTCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$NESTCOIND_PIDFILE" \
    --chuid $NESTCOIND_USER:$NESTCOIND_GROUP \
    --exec "$NESTCOIND_BIN" \
    -- \
    -pid="$NESTCOIND_PIDFILE" \
    -conf="$NESTCOIND_CONFIGFILE" \
    -datadir="$NESTCOIND_DATADIR" \
    -disablewallet \
    -daemon

